# dotfiles/etc/systemd/system/wol.service
[Unit]
Description=Wake-on-LAN

[Service]
Type=oneshot
ExecStart=/sbin/ethtool -s ADAPTER_NAME wol g

[Install]
WantedBy=basic.target
